# etc/profile.d/welcome.sh
echo -e '\033[01;33mWelcome to Zirafa!                   '
echo -e '\033[01;33m                                     '
echo -e '\033[01;33mOwner: Michal Hozza                  '
echo -e '\033[01;33mEmail: mhozza@gmail.com              '
echo -e '\033[01;33m                                     '
echo -e '\033[01;33m              "9Qk.  )U8e`           '
echo -e '\033[01;33m               :3ma%_iXs.            '
echo -e '\033[01;33m             ._a,,   _aa,.           '
echo -e '\033[01;33m           __aaawmWE` vmmmmmBmQQQmc. '
echo -e '\033[01;33m        _ayQQQQDTYSs._mQWWWQQQWWQQe  '
echo -e '\033[01;33m        3VBHT!~   ~9mmQE+~"???????^  '
echo -e '\033[01;33m                   -3WWQg%           '
echo -e '\033[01;33m                     3UWW#(          '
echo -e '\033[01;33m                     --+X-+.         '
echo -e '\033[01;33m _a,,                 .jm,wm;        '
echo -e '\033[01;33m 3QQz                 -*$QQQo.       '
echo -e '\033[01;33m :V8(                  .)QQ@$o.      '
echo -e '\033[01;33m <o=:                 .<wQQ()Z;      '
echo -e '\033[01;33m -mm;                .vQWHQmmWp.     '
echo -e '\033[01;33m _mQw,.       __aawym+]Qf dBWWk`     '
echo -e '\033[01;33m -dQrWQ()WWQQQLmQWBWQgmQp%d>)&`      '
echo -e '\033[01;33m  ~VmWWpwW("QW#QQW~]Q@Y$QWQmd}       '
echo -e '\033[01;33m   -"HWYWWgyW+)(4Q;jWj =#^+QWz       '
echo -e '\033[01;33m     :v,dQQQQaw>)??$Wo%wm,=QW=.      '
echo -e '\033[01;33m     .3mQQQQBVQQo, )VQQWWXmQWz`      '
echo -e '\033[01;33m      |WB!Y"(:QWF> --3WV`-??QL_>     '
echo -e '\033[01;33m     <mWk_> =qWQp% <wmk-. =)QQQm,    '
echo -e '\033[01;33m     ]QQQQp .VWWQk )QQmwa.-?VVVVS.   '
echo -e '\033[01;33m                                     '
